# host/jobs.cron
* * * * * %(current)s/bin/python %(current)s/manage.py send_emails
